# images/files/inspec.service
[Unit]
Description=run inspec tests

[Service]
User=root
Type=oneshot
ExecStart=/opt/chef/inspec/run_inspec.sh